# xgmobile-manager/assets/services/egpu-boot-guardian.service
[Unit]
Description=XG Mobile Boot Safety Check
Before=sddm.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/egpu-boot-guardian
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
